# data_source/env_mas/configs/crontab/crontab.conf
*/10 * * * * "/data/configs/crontab/clean.sh"
